# roles/vhosts/acp_server_gemini/templates/gemini.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>plus.svc.xworkspace.acp.gemini</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            export PATH="{{ acp_gemini_path }}"
            export LITELLM_MASTER_KEY="{{ acp_gemini_auth_token }}"
            {% for key, value in acp_gemini_environment.items() %}
            export {{ key }}="{{ value }}"
            {% endfor %}
            
            exec "{{ acp_gemini_bridge_binary_path }}" acp-server \
                --port {{ acp_gemini_listen_port }} \
                --host {{ acp_gemini_listen_host }} \
                --sub-command gemini \
                --sub-command mcp-app-server
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>{{ acp_gemini_workdir }}</string>
    <key>StandardOutPath</key>
    <string>{{ ansible_env.HOME }}/.local/state/xworkspace/acp.gemini.log</string>
    <key>StandardErrorPath</key>
    <string>{{ ansible_env.HOME }}/.local/state/xworkspace/acp.gemini.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{ acp_gemini_workdir }}</string>
    </dict>
</dict>
</plist>
