# roles/vhosts/xray-exporter/templates/xray-exporter.service.j2
[Unit]
Description=Xray Exporter (Binary)
After=network.target

[Service]
Type=simple
User={{ xray_exporter_user }}
Group={{ xray_exporter_group }}
Environment="XRAY_STATS_URL={{ xray_exporter_stats_url }}"
Environment="LISTEN_ADDR={{ xray_exporter_listen_addr }}"
ExecStart={{ xray_exporter_binary_path }}
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
