# roles/vhosts/gateway_openclaw/templates/openclaw.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>plus.svc.xworkspace.openclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            export PATH="{{ gateway_openclaw_service_path }}"
            export NODE_COMPILE_CACHE="{{ gateway_openclaw_compile_cache_dir }}"
            export OPENCLAW_NO_RESPAWN=1
            
            exec "{{ gateway_openclaw_binary_path }}" gateway --json-config "{{ gateway_openclaw_config_path }}"
        </string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>{{ gateway_openclaw_home }}</string>
    <key>StandardOutPath</key>
    <string>{{ gateway_openclaw_home }}/.local/state/xworkspace/openclaw.log</string>
    <key>StandardErrorPath</key>
    <string>{{ gateway_openclaw_home }}/.local/state/xworkspace/openclaw.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>{{ gateway_openclaw_home }}</string>
    </dict>
</dict>
</plist>
