# roles/vhosts/action-runner/templates/github-runner.service.j2
[Unit]
Description=GitHub Actions Runner
After=network.target

[Service]
ExecStart=/home/runner/actions-runner/run.sh
WorkingDirectory=/home/runner/actions-runner
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
